# root/app/openvpn/supervisord.template.conf
[program:openvpn-{VPN_COUNTRY}]
command = /usr/sbin/openvpn --script-security 2 --config /app/openvpn/config-{VPN_COUNTRY}.ovpn
autorestart = true
priority = 1
startsecs = 1
stdout_logfile=/var/log/openvpn-{VPN_COUNTRY}.log
logfile_maxbytes = 50MB
